# File Monitoring Tool/file-monitor.service
[Unit]
Description=File Monitoring Service
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/file_monitor.py --verbose
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
